# server/aesdsocket-start-stop
#!/bin/sh
### BEGIN INIT INFO
# Provides:          aesdsocket
# Required-Start:    $network $local_fs
# Required-Stop:     $network $local_fs
# Default-Start:     S
# Default-Stop:      0 6
# Short-Description: Start/stop aesdsocket daemon
### END INIT INFO

DAEMON=/usr/bin/aesdsocket
PIDFILE=/var/run/aesdsocket.pid
DAEMON_OPTS="-d"

case "$1" in
  start)
    echo "Starting aesdsocket..."
    if command -v start-stop-daemon >/dev/null 2>&1; then
      start-stop-daemon -S -x $DAEMON -- $DAEMON_OPTS
    else
      nohup $DAEMON $DAEMON_OPTS >/dev/null 2>&1 &
      echo $! > "$PIDFILE" || true
    fi
    ;;
  stop)
    echo "Stopping aesdsocket..."
    if [ -f "$PIDFILE" ]; then
      kill -TERM $(cat "$PIDFILE") 2>/dev/null || true
      rm -f "$PIDFILE"
    else
      pkill -f $(basename $DAEMON) || true
    fi
    ;;
  restart)
    $0 stop
    sleep 1
    $0 start
    ;;
  status)
    if pidof $(basename $DAEMON) >/dev/null 2>&1; then
      echo "aesdsocket running"
      exit 0
    else
      echo "aesdsocket not running"
      exit 1
    fi
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 2
    ;;
esac

exit 0
