# server/aesdsocket-start-stop
#!/bin/sh
### BEGIN INIT INFO
# Provides:          aesdsocket
# Required-Start:    $network $syslog
# Required-Stop:     $network $syslog
# Default-Start:     S
# Default-Stop:      0 6
# Short-Description: Start/Stop aesdsocket daemon
### END INIT INFO

DAEMON=/usr/bin/aesdsocket
DAEMON_OPTS="-d"
PIDFILE=/var/run/aesdsocket.pid
LOGTAG="aesdsocket"

case "$1" in
  start)
    echo "Starting aesdsocket..."
    # start-stop-daemon will create the pidfile; do NOT use --background because the daemon will daemonize (-d)
    start-stop-daemon --start --oknodo --make-pidfile --pidfile "${PIDFILE}" --exec "${DAEMON}" -- ${DAEMON_OPTS} || {
      echo "Failed to start aesdsocket" >&2
      exit 1
    }
    ;;
  stop)
    echo "Stopping aesdsocket..."
    start-stop-daemon --stop --pidfile "${PIDFILE}" --retry 5 || {
      echo "Failed to stop aesdsocket" >&2
      exit 1
    }
    # ensure pidfile removed
    [ -f "${PIDFILE}" ] && rm -f "${PIDFILE}"
    ;;
  restart)
    $0 stop
    sleep 1
    $0 start
    ;;
  status)
    if [ -f "${PIDFILE}" ] && kill -0 $(cat "${PIDFILE}") 2>/dev/null; then
      echo "aesdsocket running (pid $(cat ${PIDFILE}))"
      exit 0
    else
      echo "aesdsocket not running"
      exit 1
    fi
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 2
    ;;
esac

exit 0

